# systemd/user/opencode-upstream.service
[Unit]
Description=OpenCode Upstream Server (user)
After=network.target

[Service]
Type=simple
Environment=PATH=%h/micromamba/envs/main/bin:%h/.local/share/mamba/envs/main/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/fish -c 'if test -f "$HOME/.config/opencode/server-env.fish"; source "$HOME/.config/opencode/server-env.fish"; end; set -q OPENCODE_HOSTNAME; or set -gx OPENCODE_HOSTNAME 0.0.0.0; set -q OPENCODE_PORT; or set -gx OPENCODE_PORT 14096; exec opencode serve --hostname "$OPENCODE_HOSTNAME" --port "$OPENCODE_PORT"'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
